# configs/openrc/wgwd
#!/sbin/openrc-run
# https://github.com/OpenRC/openrc/blob/master/service-script-guide.md

command="/usr/local/bin/wgwd"
command_args="get -remote_interface=pppoe0 -wg_interface=wg0 -interval=5m webdav -endpoint=https://192.168.1.2/ -username=admin -password=adminadmin -filepath=/dav/center.json -encryption_key=adminadmin"
command_background=true
command_user="root"
pidfile="/run/${RC_SVCNAME}.pid"

depend() {
  need net
  use dns logger netmount
  after sshd firewall
}
