# configs/openrc/wgwd_s3
#!/sbin/openrc-run
# https://github.com/OpenRC/openrc/blob/master/service-script-guide.md
# env
export W_REMOTE_INTERFACE="xxx"
export W_WG_INTERFACE="xxx"
export W_ENDPOINT="xxx"
export W_ACCESS_KEY_ID="xxx"
export W_SECRET_ACCESS_KEY="xxx"
export W_BUCKET="xxx"
export W_OBJECT_PATH="xxx"
export W_ENCRYPTION_KEY="xxx"
export N_INTERVAL="5m"

command="/usr/local/bin/wgwd"
command_args="get -remote_interface=${W_REMOTE_INTERFACE} -wg_interface=${W_WG_INTERFACE} -interval=${N_INTERVAL} s3 -endpoint=${W_ENDPOINT} -path_style -access_key_id=${W_ACCESS_KEY_ID} -secret_access_key=${W_SECRET_ACCESS_KEY} -bucket=${W_BUCKET} -object_path=${W_OBJECT_PATH} -encryption_key=${W_ENCRYPTION_KEY}"
command_background=true
command_user="root"
pidfile="/run/${RC_SVCNAME}.pid"

depend() {
  need net
  use dns logger netmount
  after firewall wireguard
}
