# configs/nwg.initd
#!/bin/sh /etc/rc.common
START=90
USE_PROCD=1
COMMAND="/usr/bin/nwg"
Command_Args="-id DEVICE_UUID -url NETINFO_RECORD_URL -interface pppoe0 -wg_interface wg0 -interval 30m"

start_service() {
  procd_open_instance
  procd_set_param respawn
  procd_set_param command "$COMMAND" $Command_Args
  procd_close_instance
}