# configs/wgwd.openrc
#!/sbin/openrc-run
# https://github.com/OpenRC/openrc/blob/master/service-script-guide.md

command="/usr/local/bin/nwg"
command_args="get -remote_interface=pppoe0 -wg_interface=wg0 nconnect -id=center -endpoint=http://localhost:1996/ -interval=5s"
command_background=true
command_user="root"
pidfile="/run/${RC_SVCNAME}.pid"

depend() {
  need net
  use dns logger netmount
  after sshd firewall
}
