# configs/wgwd.initd
#!/bin/sh /etc/rc.common
START=90
USE_PROCD=1
COMMAND="/usr/bin/wgwd"
Command_Args="get -remote_interface=pppoe0 -wg_interface=wg0 nconnect -id=center -endpoint=http://localhost:1996/ -interval=5s"

start_service() {
  procd_open_instance
  procd_set_param respawn
  procd_set_param command "$COMMAND" $Command_Args
  procd_close_instance
}